# net-misc/broadlink-mqtt/files/broadlink-mqtt
#!/sbin/openrc-run

depend() {
	need net
}

checkconfig() {
        if [ ! -f /etc/${SVCNAME}.conf ] ; then
                eerror "No /etc/${SVCNAME}.conf file exists!"
        fi
}

start() {
	checkconfig || return 1
        ebegin "Starting ${SVCNAME}"
        start-stop-daemon --start --exec /usr/bin/${SVCNAME}.py --background --make-pidfile --pidfile /var/run/${SVCNAME}.pid
        eend $?
}

stop() {
        ebegin "Stopping ${SVCNAME}"
        start-stop-daemon --stop --exec /usr/sbin/${SVCNAME}.py --pidfile /var/run/${SVCNAME}.pid
        eend $?
}

